# pods/init.d/init
#!/sbin/openrc-run


name=$RC_SVCNAME
description="Podman Pods Service"

command="/usr/bin/podman"
command_args="kube play --network ${pod_network} /home/${command_user}/pods/${pod_name}/pod.yaml >dev/null  2>&1 &"


depend() {
	need net # waits for net to start
}

cleanup() {
	/sbin/runuser -u ${command_user} ${command} pod exists ${pod_name}
	result=$?
	if [ $result -eq 0 ]; then
		/sbin/runuser -u ${command_user} ${command} pod stop ${pod_name}
		/sbin/runuser -u ${command_user} ${command} pod rm ${pod_name}
	fi
}

start_pre() {
	cleanup
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
	cleanup
	eend $?
}

status() {
	/sbin/runuser -u {$command_user} -- ${command} pod exists ${pod_name} 2>/dev/null
	result=$?
	if [ $result -eq 0 ]; then
		einfo "status: started"
		return 0
	else
		einfo "status: stopped"
		return 3
	fi
}
