# pods/init.d/init
#!/sbin/openrc-run


name=$RC_SVCNAME
description="Podman Pods Service"

command="/usr/bin/podman"
command_args="kube play --replace /home/${command_user}/pods/${pod_name}/pod.yaml > /dev/null && /sbin/runuser -u ${command_user} ${command} pod start ${pod_name} > /dev/null " 

depend() {
	need net # waits for net to start
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
        /sbin/runuser -u ${command_user} ${command} pod stop ${pod_name} > /dev/null
#	/sbin/runuser -u ${command_user} ${command} pod rm ${pod_name} > /dev/null
	eend $?
}

#start() {
#	ebegin "Starting $RC_SVCNAME"
#	cd /home/${command_user}/pods/${pod_name}
#	/sbin/runuser -u ${command_user} -- podman-compose up --force-recreate -d #> /dev/null
#	eend $?
#}
 
status() {
	/sbin/runuser -u {$command_user} -- ${command} pod exists ${pod_name} 2>/dev/null
	result=$?
	if [ $result -eq 0 ]; then
		einfo "status: started"
		return 0
	else
		einfo "status: stopped"
		return 3
	fi
}
