# pods/init.d/pod-michel
#!/sbin/openrc-run
# this file has a symlink in /etc/init.d

pod_name=michel
command_user=liam
pod_network=michel

source "/home/${command_user}/pods/init.d/init"
